# docker/supervisor.conf
[program:messenger]
user=root
directory=/root/messenger/backend
command=python3 manage.py runserver 0.0.0.0:8000
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0